# debian/cron.d
PATH=/bin:/usr/bin:/sbin:/usr/sbin

* * * * * root ganeti-failover-gateway-cluster /var/lib/ganeti-failover-gateway-cluster/gw.txt
